# docker/images/msl-api/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisord.log
pidfile=/var/run/supervisord.pid

[program:nginx]
command=/usr/sbin/nginx -g 'daemon off;'
autostart=true
autorestart=true
stdout_logfile=/var/log/nginx-supervisord.log
stderr_logfile=/var/log/nginx-supervisord-error.log
user=root

[program:laravel-queue]
command=/usr/bin/php8.0 /var/www/msl_api/artisan queue:work  --rest=1 --tries=3 --timeout=300
autostart=true
autorestart=true
stdout_logfile=/var/log/laravel-queue-supervisord.log
stderr_logfile=/var/log/laravel-queue-supervisord-error.log
directory=/var/www/msl_api
user=www-data

[group:laravel]
programs=nginx,laravel-queue
