# docker/images/msl-api/supervisord-webserver-only.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisord.log
pidfile=/var/run/supervisord.pid

[program:nginx]
command=/usr/sbin/nginx -g 'daemon off;'
autostart=true
autorestart=true
stdout_logfile=/var/log/nginx-supervisord.log
stderr_logfile=/var/log/nginx-supervisord-error.log
user=root

[program:php-fpm]
command=/usr/sbin/php-fpm8.3 --nodaemonize --fpm-config /etc/php/8.3/fpm/php-fpm.conf
autostart=true
autorestart=true
stdout_logfile=/var/log/php-fpm-supervisord.log
stderr_logfile=/var/log/php-fpm-supervisord-error.log
user=root

[group:msl-api]
programs=nginx,php-fpm
